# tools/local/systemd/saltmarcher-update.timer
[Unit]
Description=Run SaltMarcher stable updater

[Timer]
OnCalendar=07:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
