# tools/local/systemd/saltmarcher-update.service
[Unit]
Description=SaltMarcher stable updater

[Service]
Type=oneshot
ExecStart=%h/.local/bin/saltmarcher-update.sh
